# guest-API-SDK/17.06/wrs-guest-scale-2.0.4/scripts/guest-scale-agent.service
[Unit]
Description=Guest Scale Agent
After=cloud-init.service

[Service]
ExecStart=/usr/sbin/guest_scale_agent
Type=simple
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
